# deb-packages/csa_0.0.1_all/DEBIAN/prerm
#!/bin/sh
if [ "$1" = "remove" ]; then
	update-rc.d -f csa remove
fi
