# deb-packages/csm_0.0.1_all/DEBIAN/postinst
#!/bin/sh
update-rc.d -f csm defaults
rabbitmq-plugins enable rabbitmq_management
service rabbitmq-server restart
